# docs/examples/vm-builder-agent.service
[Unit]
Description=vm-builder-agent
Documentation=https://github.com/tlhakhan/vm-builder-agent
After=network-online.target libvirtd.service
Wants=network-online.target
Requires=libvirtd.service

[Service]
ExecStartPre=/bin/mkdir -p /var/lib/vm-builder-agent/workspaces
ExecStartPre=/bin/mkdir -p /var/lib/vm-builder-agent/cloud-image-cache
ExecStartPre=/bin/mkdir -p /etc/vm-builder-agent/private
ExecStart=/usr/local/bin/vm-builder-agent \
  --listen                :8443 \
  --agent-mtls \
  --agent-trusted-ca-url  https://pki.example.internal/vm-builder-ca.crt \
  --private-dir           /etc/vm-builder-agent/private \
  --core-repo             https://github.com/tlhakhan/vm-builder-core \
  --terraform             tofu \
  --workspaces-dir        /var/lib/vm-builder-agent/workspaces \
  --cloud-image-cache-dir /var/lib/vm-builder-agent/cloud-image-cache

Restart=on-failure
RestartSec=5s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=vm-builder-agent

# Run as root so virsh and terraform can manage KVM resources.
# Scope down once a dedicated libvirt group or polkit policy is in place.
User=root
Group=root

[Install]
WantedBy=multi-user.target
